# Template/COMMON/delete-files.list
boot/*.old-dkms
etc/*.dpkg-dist
etc/*.dpkg-old
etc/apt/sources.list.d/*.dpkg-dist
etc/apt/sources.list.d/*.dpkg-old
etc/desktop-session
#etc/grub.d/*.dpkg-dist
etc/group-
etc/gshadow-
#etc/init.d/*.dpkg-dist
etc/issue
etc/issue.net
etc/machine-id
etc/passwd-
etc/shadow-
etc/skel/.streamtuner
etc/skel/.xinitrc
etc/skel/.xinitrc-custom
etc/slim.conf.dpkg-*
etc/subgid-
etc/subuid-
#etc/systemd
etc/xdg/menus/debian-menu.menu
#lib/systemd
root/.bash_history
tmp/*
#usr/lib/x86_64-linux-gnu/systemd
usr/share/antiX/lib
#usr/share/applications/gparted.desktop
#usr/share/applications/gufw.desktop
usr/share/applications/luckybackup-kde-su.desktop
usr/share/applications/ndisgtk-kde.desktop
#usr/share/applications/synaptic-kde.desktop
#usr/share/applications/synaptic.desktop
usr/share/applications/vim.desktop
usr/share/desktop-defaults
#usr/share/fluxbox/styles/Artwiz
#usr/share/fluxbox/styles/Bity_fluxbox
#usr/share/fluxbox/styles/BlueNight
#usr/share/fluxbox/styles/Doty
#usr/share/fluxbox/styles/Emerge
#usr/share/fluxbox/styles/Flux
#usr/share/fluxbox/styles/LemonSpace
#usr/share/fluxbox/styles/Makro
#usr/share/fluxbox/styles/MerleyKay
#usr/share/fluxbox/styles/Meta
#usr/share/fluxbox/styles/Nyz
#usr/share/fluxbox/styles/Operation
#usr/share/fluxbox/styles/Outcomes
#usr/share/fluxbox/styles/Results
#usr/share/fluxbox/styles/Shade
#usr/share/fluxbox/styles/Squared*
#usr/share/fluxbox/styles/Twice
#usr/share/fluxbox/styles/arch
#usr/share/fluxbox/styles/bloe
#usr/share/fluxbox/styles/bora_black
#usr/share/fluxbox/styles/bora_blue
#usr/share/fluxbox/styles/bora_green
#usr/share/fluxbox/styles/carp
#usr/share/fluxbox/styles/debian-blue
#usr/share/fluxbox/styles/debian-dark
#usr/share/fluxbox/styles/green_tea
#usr/share/fluxbox/styles/ostrich
#usr/share/fluxbox/styles/qnx-photon
#usr/share/fluxbox/styles/zimek*
#usr/share/icewm/themes/Infadel2
#usr/share/icewm/themes/gtk
#usr/share/icewm/themes/gtk2
#usr/share/icewm/themes/icedesert
#usr/share/icewm/themes/metal2
#usr/share/icewm/themes/motif
#usr/share/icewm/themes/nice
#usr/share/icewm/themes/nice2
#usr/share/icewm/themes/warp3
#usr/share/icewm/themes/warp4
#usr/share/icewm/themes/win95
#usr/share/icewm/themes/yellowmotif
#usr/share/icons/*/icon-theme.cache
#usr/share/icons/BRIT-ICONS/icon-theme.cache
#usr/share/wallpapers/*
var/cache/apt-show-versions/*
var/cache/apt/*.bin
var/cache/debconf/*-old
var/cache/flashplugin-nonfree/flashplayer*
var/cache/flashplugin-nonfree/install*
var/cache/man/*
var/lib/apt/lists/*
var/lib/apt/periodic/*
var/lib/dbus/machine-id
var/lib/dhcp/*
var/lib/dpkg/*-old
#var/lib/systemd
var/log/*.log
var/log/*log
var/log/*mp
var/log/apt/*.log
var/log/dmesg
