# Template/COMMON/initrd/etc/udev/rules.d/90-fstab-automount.rules
ACTION=="add|remove", SUBSYSTEM=="block|mmc", KERNEL=="[hs]d[a-z]*|mmcblk[0-9]*", \
RUN+="/bin/sh -c '/sbin/make-fstab --update'"
